# board/jebrish/rg353x_minxgraph/rootfs-overlay/etc/profile.d/weston.sh
export XDG_RUNTIME_DIR=/var/weston
export WAYLAND_DISPLAY=wayland-1
export DISPLAY=:0
